# elan-center-connector.elan-rule-fetcher.cron.d
24 *	* * *	root	PYTHONPATH=/elan-agent/lib/python /elan-agent/bin/rule-fetcher
